# workloads/node/bootstrap/files/_systemd_units/cluster/setup/build-step-ca-bootstrap.service
[Unit]
Description=Build smallstep bootstrap image
Requires=k8s-network-configured.target
After=k8s-network-configured.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=manifest apply step-ca-bootstrap
ExecStartPost=resource-ready --wait step-ca-bootstrap
