# bootstrap/root/etc/systemd/system/setup-workload@phxc.service.d/deps.conf
[Unit]
Requires=workload@smallstep.target
After=workload@smallstep.target
