# bootstrap/root/etc/systemd/system/setup-workload@phxc.service.d/deps.conf
[Unit]
Requires=k3s.target
After=k3s.target
