# bootstrap/root/etc/systemd/system/setup-workload@image-registry.service.d/deps.conf
[Unit]
Requires=workload@cert-manager-issuers.target workload@longhorn.target workload@smallstep.target
After=workload@cert-manager-issuers.target workload@longhorn.target workload@smallstep.target
