# bootstrap/root/etc/systemd/system/setup-workload@image-registry.service.d/deps.conf
[Unit]
Requires=workload@longhorn.target workload@smallstep.target workload@phxc.target
After=workload@longhorn.target workload@smallstep.target workload@phxc.target
