# bootstrap/systemd-units/initramfs/dropbear-run-dir.service
[Unit]
Description=Creates dropbear dir in /run
DefaultDependencies=no

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=mkdir /run/dropbear
