# workloads/node/bootstrap/assets/_systemd_units/cluster/setup/k8s-network-configured.target
[Unit]
Description=The kubernetes network has been fully configured
Requires=configure-cilium.service apply-network-policies.service setup-coredns.service
After=configure-cilium.service apply-network-policies.service setup-coredns.service
