# bootstrap/root/etc/systemd/system/setup-workload@cert-manager-issuers.service.d/deps.conf
[Unit]
Requires=workload@cert-manager.target
After=workload@cert-manager.target
