# bootstrap/root/etc/systemd/system/setup-workload@cilium.service.d/deps.conf
[Unit]
Requires=install-cilium.service workload-ready@cilium-install.service
After=install-cilium.service workload-ready@cilium-install.service
