# workloads/node/bootstrap/assets/usr/lib/dracut/modules.d/99phoenix-cluster/system/overlay-rw.mount
[Unit]
Description=Create the overlay tmpfs

[Mount]
What=tmpfs
Where=/overlay/rw
Options=mode=0755
Type=tmpfs
