# bootstrap/root/etc/systemd/system/setup-workload@bootstrap.service.d/deps.conf
[Unit]
Requires=workload@phxc.target workload@image-registry.target workload@apt-cacher-ng.target
After=workload@phxc.target workload@image-registry.target workload@apt-cacher-ng.target
