# workloads/node/bootstrap/assets/_systemd_units/10-data-partition/data-partition.target
[Unit]
Description=The data partition is available for mounting
Before=var-lib-phxc.mount

[Install]
WantedBy=var-lib-phxc.mount
