# bootstrap/systemd-units/cluster/k3s/import-container-images.path
[Unit]
Description=Watch for new container images
DefaultDependencies=no
Requires=create-data-dir@container-images.service
After=create-data-dir@container-images.service

[Install]
WantedBy=k3s.target

[Path]
PathExistsGlob=/var/lib/phxc/container-images/*.tar
