# workloads/bootstrap/containers/node/assets/etc/systemd/system/import-images.path
[Unit]
Description=Watch for new images to import to Spegel

[Install]
WantedBy=k3s@server.target k3s@agent.target

[Path]
PathExistsGlob=/var/lib/home-cluster/images/*.tar
