# workloads/bootstrap/containers/node/assets/etc/systemd/system/k3s.service
[Unit]
Description=Starts either k3s server or agent depending on the labels in node settings

[Install]
WantedBy=multi-user.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/k3s-start
