# workloads/bootstrap/containers/node/assets/etc/systemd/system/update-root.path
[Unit]
Description=Watch for new image matching the system architecture

[Install]
WantedBy=k3s@server.target k3s@agent.target

[Path]
PathExists=/var/lib/persistent/host-images/uefi/${ARCH}.raw
