# workloads/node/bootstrap/assets/etc/systemd/system/authn-key-available.target
[Unit]
Description=The authentication key is available
AssertPathExists=/run/initramfs/authn.key
AssertPathExists=/run/initramfs/authn.pub
