# workloads/node/bootstrap/assets/etc/systemd/system/boot-partition.target
[Unit]
Description=The boot partition has been created
Before=boot.mount

[Install]
WantedBy=boot.mount
