# workloads/node/bootstrap/assets/etc/systemd/system/collect-node-state.service
[Unit]
Description=Collect system details and write them to the node state

[Install]
WantedBy=multi-user.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/lib/upkg/.upkg/home-cluster/workloads/node/bin/node-state collect
