# workloads/node/bootstrap/assets/etc/systemd/system/create-persistent-dir@.service
[Unit]
Description=Create data directory for bind mounts
RequiresMountsFor=/var/lib/persistent

[Service]
ExecStart=mkdir -p /var/lib/persistent/%i
