# workloads/node/bootstrap/assets/etc/systemd/system/decrypt-node-config.service
[Unit]
Description=Decrypt the node-config
DefaultDependencies=no
Requires=keys-available.target
Wants=copy-node-config.service download-node-config.service
After=copy-node-config.service download-node-config.service keys-available.target
Before=node-config-available.target
AssertPathExists=/run/initramfs/node-config.enc.json

[Install]
WantedBy=node-config-available.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=decrypt-node-config /run/initramfs/node-config.enc.json /run/initramfs/node-config.json
