# workloads/node/bootstrap/assets/etc/systemd/system/encryption-key-available.target
[Unit]
Description=The encryption key is available
AssertPathExists=/var/lib/systemd/credential.secret
