# workloads/node/bootstrap/assets/etc/systemd/system/import-container-images.path
[Unit]
Description=Watch for new container images

[Install]
WantedBy=k3s@server.service k3s@agent.service

[Path]
PathExistsGlob=/var/lib/persistent/container-images/*.tar
