# workloads/node/bootstrap/assets/etc/systemd/system/import-container-images.service
[Unit]
Description=Import container images to containerd

[Service]
Type=oneshot
ExecStart=import-container-images
