# workloads/node/bootstrap/assets/etc/systemd/system/k3s@.target
[Unit]
Description=The startup target for k3s server or agent
Wants=k3s@%i.service
