# workloads/node/bootstrap/assets/etc/systemd/system/node-config.target
[Unit]
Description=The node-config is available
AssertPathExists=/run/initramfs/node-config.json

[Install]
WantedBy=sysinit.target
