# workloads/node/bootstrap/assets/etc/systemd/system/report-node-state.path
[Unit]
Description=Watch for node state changes

[Install]
WantedBy=paths.target

[Path]
PathModified=/run/initramfs/node-state.json
