# workloads/node/bootstrap/assets/etc/systemd/system/report-initial-node-state.service
[Unit]
Description=Report the initial node-state to the boot-server
DefaultDependencies=no
Requires=systemd-timesyncd.service authn-key-available.target
After=systemd-timesyncd.service authn-key-available.target

[Install]
WantedBy=sysinit.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=report-node-state
