# workloads/node/bootstrap/assets/etc/systemd/system/resource-ready@.service
[Unit]
Description=Waits until %i is ready

[Service]
Type=oneshot
ExecStart=resource-ready --wait %i
