# workloads/node/bootstrap/assets/etc/systemd/system/update-node-config.service
[Unit]
Description=Download the node configuration

[Install]
WantedBy=multi-user.target

[Service]
Type=oneshot
ExecStart=update-node-config
