# workloads/node/bootstrap/assets/etc/systemd/system/workload-ready@.target
[Unit]
Description=Reached when %i is ready
Wants=workload-ready@%i.service
After=workload-ready@%i.service
