# bootstrap/systemd-units/networking/configure-networks.service
[Unit]
Description=Configure systemd-networks using node-config
DefaultDependencies=no
Before=network-pre.target
ConditionPathExists=/efi/phxc/node.json

[Install]
WantedBy=sysinit.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecCondition=get-config node fixedIPs
ExecStart=bash -eo pipefail
StandardInputText=for macaddr in $(jq -r '(.fixedIPs // {}) | keys[]' /efi/phxc/node.json); do \
    jq -r --arg macaddr "$macaddr" '\
      ([.fixedIPs[$macaddr][] | "Address=\\(.)"] | join("\\n")) as $addrs | \
      "[Match]\\nMACAddress=\\($macaddr)\\n[Network]\\n\\($addrs)\\n"' \
      /efi/phxc/node.json> /etc/systemd/network/${macaddr//:/-}.network \
  done
