# bootstrap/root/etc/systemd/system/setup-workload@coredns.service.d/deps.conf
[Unit]
Requires=workload@cilium.target workload@network-policies.target
After=workload@cilium.target workload@network-policies.target
