# bootstrap/root/etc/systemd/system/setup-workload@apt-cacher-ng.service.d/deps.conf
[Unit]
Requires=workload@phxc.target
After=workload@phxc.target
