# workloads/node/bootstrap/systemd-units/cluster/setup/k3s-network-ready.target
[Unit]
Description=The cluster network is fully set up
Requires=k3s@server.service workload@cilium.target workload@network-policies.target workload@coredns.target
After=k3s@server.service workload@cilium.target workload@network-policies.target workload@coredns.target
