# workloads/node/bootstrap/systemd-units/cluster/setup/workload@.target
[Unit]
Description=Applies the %i manifest and waits for it to be ready
Requires=setup-workload@%i.service workload-ready@%i.service
After=setup-workload@%i.service workload-ready@%i.service
