# workloads/node/bootstrap/assets/_systemd_units/25-persistent-partition/data-partition.target
[Unit]
Description=The data partition has been created
Before=var-lib-phoenix\x2dcluster.mount

[Install]
WantedBy=var-lib-phoenix\x2dcluster.mount
