# workloads/node/bootstrap/assets/_systemd_units/25-persistent-partition/persistent-partition.target
[Unit]
Description=The persistent partition has been created
Before=var-lib-persistent.mount

[Install]
WantedBy=var-lib-persistent.mount
