# bootstrap/systemd-units/update-boot/init-efi-bootmenu.service
[Unit]
Description=Populate the EFI boot menu with Phoenix Cluster entries
Requires=boot.mount
After=boot.mount

[Install]
WantedBy=multi-user.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=efi-bootentry add /boot/EFI/BOOT/BOOT${EFI_ARCH}.EFI "Phoenix Cluster"
ExecStart=efi-bootentry add /boot/EFI/BOOT/BOOT.TRY.EFI "Phoenix Cluster (new image)"
ExecStart=efi-bootentry deactivate /boot/EFI/BOOT/BOOT.TRY.EFI
