# workloads/node/bootstrap/assets/usr/lib/dracut/modules.d/99home-cluster/units/create-node-state@.service
[Unit]
Description=Create /run/initramfs/node-state.json

[Install]
WantedBy=sysinit.target

[Service]
Type=oneshot
ExecStart=/usr/bin/create-node-state %i
