# workloads/node/bootstrap/assets/_systemd_units/cluster/k3s/k3s.target
[Unit]
Description=The startup target for k3s server and agent
Wants=k3s@server.service k3s@agent.service

[Install]
WantedBy=multi-user.target
