# bootstrap/root/etc/systemd/system/setup-workload@network-policies.service.d/deps.conf
[Unit]
Requires=workload@cilium.target
After=workload@cilium.target
