# workloads/node/bootstrap/assets/_systemd_units/05-keys/surrogate-key-available.target
[Unit]
Description=The surrogate key for decrypting the systemd credential key is available
ConditionPathExists=/run/initramfs/surrogate.key.json
