# workloads/node/bootstrap/assets/_systemd_units/80-k3s/resource-ready@.target
[Unit]
Description=Reached when %i is ready
Wants=resource-ready@%i.service
After=resource-ready@%i.service
