# workloads/node/bootstrap/assets/_systemd_units/80-k3s/resource-ready@.service
[Unit]
Description=Waits until %i is ready
DefaultDependencies=no

[Service]
Type=oneshot
ExecStart=resource-ready --wait %i
