# workloads/node/bootstrap/assets/_systemd_units/80-k3s/var-lib-longhorn.mount
[Unit]
Description=/var/lib/longhorn
DefaultDependencies=no
Requires=create-persistent-dir@longhorn.service
After=create-persistent-dir@longhorn.service
Before=k3s.target

[Install]
RequiredBy=k3s.target

[Mount]
What=/var/lib/persistent/longhorn
Where=/var/lib/longhorn
Type=none
Options=bind
