# bootstrap/root/etc/profile.d/phxc-packages.sh
#!/usr/bin/env sh
PATH=/var/lib/phxc/.upkg/.bin:$PATH
