# bootstrap/root/etc/profile.d/phxc.sh
#!/usr/bin/env sh
PATH=/var/lib/phxc/.upkg/.bin:$PATH
export KUBECONFIG=/etc/rancher/k3s/k3s-node.yaml
