# workloads/node/bootstrap/assets/_systemd_units/15-boot-partition/boot.mount
[Unit]
Description=/boot

[Mount]
What=PARTUUID=${BOOT_UUID}
Where=/boot
Type=vfat
Options=
