# bootstrap/root/etc/systemd/system/setup-workload@backup-server.service.d/deps.conf
[Unit]
Requires=workload@phxc.target workload@smallstep.target trust-smallstep-root.service
After=workload@phxc.target workload@smallstep.target trust-smallstep-root.service
