# workloads/node/bootstrap/systemd-units/data-partition/enroll/disk-encryption-keys-enrolled.target
[Unit]
Description=All disk-encryption-keys have been enrolled

[Install]
WantedBy=multi-user.target
