# bootstrap/root/etc/systemd/system/setup-workload@cert-manager.service.d/deps.conf
[Unit]
Requires=k3s-network-ready.target
After=k3s-network-ready.target
