# workloads/node/bootstrap/assets/_systemd_units/90-node-config/update-node-config.timer
[Unit]
Description=Check the node-config for updates every minute

[Install]
WantedBy=timers.target

[Timer]
OnCalendar=*-*-* *:*:00
RandomizedDelaySec=59
FixedRandomDelay=true
Persistent=true
