# workloads/node/bootstrap/assets/usr/lib/dracut/modules.d/99home-cluster/system/rootimg.target
[Unit]
Description=The rootimg is ready to be mounted
